# siri/pull_data.cron
*/5 * * * * python3 fetch_and_store_arrivals.py
